# deploy/systemd/dialexis.service
[Unit]
Description=Dialexis — Reflexive Philosophy Infrastructure
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dialexis
Group=dialexis
WorkingDirectory=/opt/dialexis
Environment=DIALEXIS_CONTACT=__CONTACT__
ExecStart=/opt/dialexis/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=3
# hardening
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/opt/dialexis/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
